# bin/ssp-serv
#!/bin/sh
#
# chkconfig: - 85 15
# description:SSP Server 1.0.1
# Comments to support LSB init script conventions

# See how we were called.
case $1 in
    start|stop|restart|status)
        "{SSP_DIR}/ssp" -H -f "{SSP_DIR}/init.php" -s $1
        if [ "$1" = "stop" ]; then
	    if [ -s "{SSP_DIR}/daemon.pid" ];then
	        kill `cat {SSP_DIR}/daemon.pid`
            fi
	    exit $?
        fi
        if [ "$1" = "status" ]; then
	    exit $?
        fi
        if [ -s "{SSP_DIR}/daemon.pid" ];then
	    if [ -d "/proc/`cat {SSP_DIR}/daemon.pid`" ];then
	        exit $?
	    fi
        fi
        "{SSP_DIR}/daemon" "{SSP_DIR}/log/ssp.pid" "service {SERV_NAME} start" > "{SSP_DIR}/log/daemon.log"
        ;;
    path)
        echo "{SSP_DIR}"
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status|path}"
        exit 2
esac

exit $?