# bin/ssp-serv
#!/bin/sh
#
# chkconfig: - 85 15
# description:SSP Server 1.0.1
# Comments to support LSB init script conventions

BIN_DIR=/opt/ssp/bin
SRC_DIR={SSP_DIR}
PID_FILE=${SRC_DIR}/log/ssp.pid
DPID_FILE=${SRC_DIR}/log/ssp.pid.daemon

# See how we were called.
case $1 in
    start|stop|restart|status)
        "${BIN_DIR}/ssp" -H -f "${SRC_DIR}/init.php" --port {SERV_PORT} --pidfile "{SERV_PIDFILE}" --nthreads {SERV_NTHEADS} --max-clients {SERV_CLIENTS} --user {SERV_USER} -s $1
        if [ "$1" = "stop" ]; then
	    if [ -s "${DPID_FILE}" ];then
	        kill `cat ${DPID_FILE}`
            fi
	    exit $?
        fi
        if [ "$1" = "status" ]; then
	    exit $?
        fi
        if [ -s "${DPID_FILE}" ];then
	    if [ -d "/proc/`cat ${DPID_FILE}`" ];then
	        exit $?
	    fi
        fi
        "${BIN_DIR}/daemon" "${PID_FILE}" "service {SERV_NAME} start" >> "${SRC_DIR}/log/daemon.log"
        ;;
    path)
        echo "${SRC_DIR}"
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status|path}"
        exit 2
esac

exit $?